# rpi.service
[Unit]
Description=Python Task Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/rpi-ping/
ExecStart=python3 /home/pi/rpi-ping/main.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
Restart=always